# etc/init/lightdm-start.conf
description     "lightdm-start"
author		"Gerald Dachs <gda@dachsweb.de>"
 
start on starting lightdm

task

script
	svdrpsend PLUG softhddevice DETA
end script
